# src/monitoring/freeipa_health_textfile/templates/freeipa-health-export.timer.j2
[Unit]
Description=Run FreeIPA health metrics exporter periodically

[Timer]
OnBootSec={{ freeipa_health_initial_delay }}
OnUnitActiveSec={{ freeipa_health_interval }}
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
