# src/infrastructure/vyos_auto_update/files/vyos-updater.timer
[Unit]
Description=Run VyOS update monthly

[Timer]
OnCalendar=monthly
RandomizedDelaySec=28d
Persistent=true

[Install]
WantedBy=timers.target
