# src/infra_services/apt_cacher_ng/files/apt-cacher-ng-cleanup.timer
[Unit]
Description=Periodic cleanup of stale apt-cacher-ng InRelease temporary files

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
