# src/storage/backup_pull_sync/templates/backup-prune.timer.j2
[Unit]
# {{ ansible_managed }}
Description=Daily timer to prune old backup generations

[Timer]
OnCalendar=daily
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
