# src/network_services/cloudflared/templates/cloudflared.initd.j2
#!/sbin/openrc-run
# OpenRC init script for cloudflared
# Managed by Ansible - do not edit manually

name="cloudflared"
description="Cloudflare Tunnel daemon"

command="{{ cloudflared_binary_path }}"
command_args="tunnel run --token {{ cloudflared_token }}"
command_user="root"
supervisor="supervise-daemon"
pidfile="/run/${RC_SVCNAME}.pid"

respawn_delay={{ cloudflared_respawn_delay }}
respawn_max={{ cloudflared_respawn_max }}
respawn_period={{ cloudflared_respawn_period }}

output_log="/var/log/cloudflared/cloudflared.log"
error_log="/var/log/cloudflared/cloudflared.log"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /var/log/cloudflared
}
